# operation_system/kursWork/kwread.service
[Unit]
Description=kwread

[Service]
User=user
Group=user
ExecStart=/usr/sbin/kwreadelf

[Install]
WantedBy=multi-user.target
